# systemd/markdown.service
[Unit]
Description=Markdown static site server
After=network.target

[Service]
Type=simple
Environment=PORT=3001
WorkingDirectory=/home/andy/github/markdown
ExecStart=/home/andy/.nvm/versions/node/v24.9.0/bin/node --watch server.js
Restart=on-failure
RestartSec=5
User=andy
Group=andy

[Install]
WantedBy=multi-user.target
